# configs/lxqt/lxqt-overlay/etc/skel/.config/autostart/Connman.desktop
[Desktop Entry]
Exec=/usr/bin/runcmst
Name=Connman
Type=Application
Version=1.0
X-LXQt-Need-Tray=true
