# configs/lxqt/lxqt-overlay/etc/skel/.config/autostart/NumLock.desktop
[Desktop Entry]
Exec=/usr/bin/numlockx off
Name=NumLock
Type=Application
Version=1.0
