# configs/shared/overlay/etc/local.d/rc.local.shutdown.stop
#!/bin/bash
#
# /etc/rc.local.shutdown.stop: Local shutdown script.
#

if [ -e /etc/rc.local.shutdown ]; then
	# Execute it
	/etc/rc.local.shutdown
fi

