# configs/shared/overlay/etc/local.d/rc.local.start
#!/bin/bash
#
# /etc/rc.local.start: Local multi-user startup script.
#

if [ -e /etc/rc.local ]; then
	# Execute it
	/etc/rc.local
fi

