# configs/shared/overlay-livecd/etc/NetworkManager/dispatcher.d/10_update_pacman_config
#!/bin/sh

case "$2" in
	up)
		/etc/manjaro/mirrorcheck
	;;
	vpn-up)
		/etc/manjaro/mirrorcheck
	;;
esac
