# configs/lxqt-openrc/lxqt-overlay/etc/skel/.config/autostart/Qasmixer.desktop
[Desktop Entry]
Exec=/usr/bin/qasmixer -t
Name=Qasmixer
Type=Application
Version=1.0
X-LXQt-Need-Tray=true
